# plugins/minutes/dev.getminutes.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.getminutes.watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/silverbook/.cargo/bin/minutes</string>
        <string>watch</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/silverbook</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/silverbook</string>
        <key>PATH</key>
        <string>/Users/silverbook/.local/bin:/Users/silverbook/.cargo/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/silverbook/.minutes/logs/watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/silverbook/.minutes/logs/watcher.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>Nice</key>
    <integer>5</integer>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
